# system/mediamtx.service
[Unit]
Description=MediaMTX
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/usr/local/bin/mediamtx /etc/mediamtx/mediamtx.yml

[Install]
WantedBy=multi-user.target
